# archivist-cleanup/crontab
*/20 * * * * /usr/sbin/cleanup.sh >> /dev/stdout 2>&1
# crontab requires an empty line at the end of the file
